# glpi-cron/glpi.cron
*/2 * * * * /usr/bin/php /var/www/glpi/front/cron.php &>/dev/null
